# systemd/claude-dir-rotate.timer
[Unit]
Description=Daily rotation of Claude Code ~/.claude/ directory

[Timer]
OnCalendar=*-*-* 23:55:00
RandomizedDelaySec=0
Persistent=true

[Install]
WantedBy=timers.target
