# systemd/claude-dir-rotate.timer
[Unit]
Description=Daily rotation of Claude Code ~/.claude/ directory

[Timer]
OnCalendar=daily
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
